# op/supervisord_quant.conf
; [inet_http_server]
; port=0.0.0.0:9001   ; 允许所有网络接口访问
; username=user       ; 可选：配置基本身份验证的用户名
; password=pass       ; 可选：配置基本身份验证的密码

[supervisord]
logfile=./logs/supervisord/supervisord.log
logfile_maxbytes=2MB
logfile_backups=100
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0777

[supervisorctl]
serverurl=unix:////tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface





[program:ETH]
command=python3 -u quantVole.py ETH-USDT
directory=./app
redirect_stderr=true
stopasgroup=true
killasgroup=true
autostart=true
autorestart=unexpected  ; 仅在非正常退出时重启
exitcodes=0             ; 正常退出的状态码为 0
stdout_logfile=./logs/supervisord/quantVole.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=100




