# op/supervisord_quant.conf
[inet_http_server]
port=0.0.0.0:9002   ; 允许所有网络接口访问
username=%(ENV_SUPERVISOR_USER)s       ; 可选：配置基本身份验证的用户名
password=%(ENV_SUPERVISOR_PASS)s       ; 可选：配置基本身份验证的密码

[supervisord]
logfile=./logs/supervisord/supervisord.log
logfile_maxbytes=2MB
logfile_backups=100
loglevel=info
pidfile=/tmp/supervisord_quant.pid
nodaemon=false

[unix_http_server]
file=/tmp/supervisor_quant.sock
chmod=0777

[supervisorctl]
serverurl=unix:////tmp/supervisor_quant.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface



[program:auto_upload_account]
command=python3 -u auto_upload_account.py
directory=./app/monitor
redirect_stderr=true
stopasgroup=true
killasgroup=true
autostart=true
autorestart=unexpected  ; 仅在非正常退出时重启
exitcodes=0             ; 正常退出的状态码为 0
stdout_logfile=./logs/supervisord/auto_upload_account.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=30


[program:quantVole-BTC]
command=python3 -u quantVole.py BTC-USDT
directory=./app
redirect_stderr=true
stopasgroup=true
killasgroup=true
autostart=false
autorestart=unexpected  ; 仅在非正常退出时重启
exitcodes=0             ; 正常退出的状态码为 0
stdout_logfile=./logs/supervisord/quantVole-BTC.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=30

[program:quantVole-ETH]
command=python3 -u quantVole.py ETH-USDT
directory=./app
redirect_stderr=true
stopasgroup=true
killasgroup=true
autostart=false
autorestart=unexpected  ; 仅在非正常退出时重启
exitcodes=0             ; 正常退出的状态码为 0
stdout_logfile=./logs/supervisord/quantVole-ETH.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=30

[program:quantVole-DOGE]
command=python3 -u quantVole.py DOGE-USDT
directory=./app
redirect_stderr=true
stopasgroup=true
killasgroup=true
autostart=false
autorestart=unexpected  ; 仅在非正常退出时重启
exitcodes=0             ; 正常退出的状态码为 0
stdout_logfile=./logs/supervisord/quantVole-DOGE.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=30

[program:quantVole-FLOKI]
command=python3 -u quantVole.py FLOKI-USDT
directory=./app
redirect_stderr=true
stopasgroup=true
killasgroup=true
autostart=false
autorestart=unexpected  ; 仅在非正常退出时重启
exitcodes=0             ; 正常退出的状态码为 0
stdout_logfile=./logs/supervisord/quantVole-FLOKI.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=30

[program:quantVole-JST]
command=python3 -u quantVole.py JST-USDT
directory=./app
redirect_stderr=true
stopasgroup=true
killasgroup=true
autostart=false
autorestart=unexpected  ; 仅在非正常退出时重启
exitcodes=0             ; 正常退出的状态码为 0
stdout_logfile=./logs/supervisord/quantVole-JST.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=30

[program:quantVole-LUNC]
command=python3 -u quantVole.py LUNC-USDT
directory=./app
redirect_stderr=true
stopasgroup=true
killasgroup=true
autostart=false
autorestart=unexpected  ; 仅在非正常退出时重启
exitcodes=0             ; 正常退出的状态码为 0
stdout_logfile=./logs/supervisord/quantVole-LUNC.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=30

[program:quantVole-ORDI]
command=python3 -u quantVole.py ORDI-USDT
directory=./app
redirect_stderr=true
stopasgroup=true
killasgroup=true
autostart=false
autorestart=unexpected  ; 仅在非正常退出时重启
exitcodes=0             ; 正常退出的状态码为 0
stdout_logfile=./logs/supervisord/quantVole-ORDI.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=30

[program:quantVole-RACA]
command=python3 -u quantVole.py RACA-USDT
directory=./app
redirect_stderr=true
stopasgroup=true
killasgroup=true
autostart=false
autorestart=unexpected  ; 仅在非正常退出时重启
exitcodes=0             ; 正常退出的状态码为 0
stdout_logfile=./logs/supervisord/quantVole-RACA.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=30

[program:quantVole-ZIL]
command=python3 -u quantVole.py ZIL-USDT
directory=./app
redirect_stderr=true
stopasgroup=true
killasgroup=true
autostart=false
autorestart=unexpected  ; 仅在非正常退出时重启
exitcodes=0             ; 正常退出的状态码为 0
stdout_logfile=./logs/supervisord/quantVole-ZIL.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=30

[program:quantVole-ZRX]
command=python3 -u quantVole.py ZRX-USDT
directory=./app
redirect_stderr=true
stopasgroup=true
killasgroup=true
autostart=false
autorestart=unexpected  ; 仅在非正常退出时重启
exitcodes=0             ; 正常退出的状态码为 0
stdout_logfile=./logs/supervisord/quantVole-ZRX.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=30




